# batch_app/cronjobs/customer.cron
0 0 * * * python /code/app/main_customer.py >> /code/log/customer.log 2>&1
